# supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat7]
command=service tomcat7 start
